# scripts/com.jingerzz.simple-portfolio-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.jingerzz.simple-portfolio-refresh</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/jxie/Simple Portfolio Manager/.venv/bin/python</string>
      <string>/Users/jxie/Simple Portfolio Manager/scripts/daily_refresh.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jxie/Simple Portfolio Manager</string>

    <key>StartCalendarInterval</key>
    <array>
      <dict>
        <key>Weekday</key><integer>1</integer>
        <key>Hour</key><integer>16</integer>
        <key>Minute</key><integer>30</integer>
      </dict>
      <dict>
        <key>Weekday</key><integer>2</integer>
        <key>Hour</key><integer>16</integer>
        <key>Minute</key><integer>30</integer>
      </dict>
      <dict>
        <key>Weekday</key><integer>3</integer>
        <key>Hour</key><integer>16</integer>
        <key>Minute</key><integer>30</integer>
      </dict>
      <dict>
        <key>Weekday</key><integer>4</integer>
        <key>Hour</key><integer>16</integer>
        <key>Minute</key><integer>30</integer>
      </dict>
      <dict>
        <key>Weekday</key><integer>5</integer>
        <key>Hour</key><integer>16</integer>
        <key>Minute</key><integer>30</integer>
      </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/jxie/Simple Portfolio Manager/daily_refresh.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jxie/Simple Portfolio Manager/daily_refresh.log</string>
  </dict>
</plist>
